# cookbooks/aws-parallelcluster-install/templates/default/base/setup-ephemeral.service.erb
[Unit]
Description=Setup ephemeral drives service
After=<%= @network_target %>

[Service]
ExecStart=/usr/local/sbin/setup-ephemeral-drives.sh

[Install]
WantedBy=multi-user.target
