# cookbooks/aws-parallelcluster-config/templates/default/base/parallelcluster_supervisord.conf.erb
# Generated by Chef for AWS ParallelCluster <%= node['cluster']['node_type'] -%>
# Local modifications could be be overwritten.
<% case node['cluster']['node_type'] -%>
<%# HeadNode -%>
<% when 'HeadNode' -%>
[program:cfn-hup]
command = bash -c "[ -f /etc/profile.d/proxy.sh ] && . /etc/profile.d/proxy.sh; <%= node['cluster']['cfn_bootstrap_virtualenv_path'] %>/bin/cfn-hup"
# The following are needed because cfn-hup starts as a daemon
exitcodes = 0
autorestart = unexpected
startsecs = 0
<% if node['cluster']['scheduler'] == 'slurm' -%>
[program:clustermgtd]
command = <%= node['cluster']['node_virtualenv_path'] %>/bin/clustermgtd
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/clustermgtd
<% end -%>
<% unless node['cluster']['scheduler'] == 'awsbatch' -%>
[program:clusterstatusmgtd]
command = <%= node['cluster']['cookbook_virtualenv_path'] %>/bin/python /opt/parallelcluster/scripts/clusterstatusmgtd.py
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/clusterstatusmgtd
<% end -%>
<% if node['cluster']['dcv_enabled'] == "head_node" && node['conditions']['dcv_supported'] -%>
[program:pcluster_dcv_authenticator]
command = <%= node['cluster']['dcv']['authenticator']['virtualenv_path'] %>/bin/python <%= node['cluster']['dcv']['authenticator']['user_home'] %>/pcluster_dcv_authenticator.py
      --port <%= Integer(node['cluster']['dcv_port']) + 1 %>
      --certificate <%= node['cluster']['dcv']['authenticator']['certificate'] %>
      --key <%= node['cluster']['dcv']['authenticator']['private_key'] %>
user = <%= node['cluster']['dcv']['authenticator']['user'] %>
environment = HOME="<%= node['cluster']['dcv']['authenticator']['user_home'] %>",USER="<%= node['cluster']['dcv']['authenticator']['user'] %>"
<% end -%>

<%# ComputeFleet -%>
<% when 'ComputeFleet' -%>
<% if node['cluster']['scheduler'] == 'slurm' -%>
[program:computemgtd]
command = <%= node['cluster']['node_virtualenv_path'] %>/bin/computemgtd
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/computemgtd
<% end -%>

<% end -%>
