# cookbooks/aws-parallelcluster-platform/templates/supervisord/parallelcluster_supervisord.conf.erb
# Generated by Chef for AWS ParallelCluster <%= node['cluster']['node_type'] -%>
# Local modifications could be overwritten.
<%# HeadNode, ComputeFleet -%>
<% case node['cluster']['node_type'] -%>
<% when 'HeadNode', 'ComputeFleet' -%>
[program:cfn-hup]
command = <%= node['cluster']['scripts_dir']%>/cfn-hup-runner.sh
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/cfn-hup-runner.log
stdout_logfile_maxbytes = 1MB
<% if @region.start_with?('us-iso') -%>
  environment = AWS_CA_BUNDLE="<%= @aws_ca_bundle %>"
<% end -%>
<% end -%>

<% case node['cluster']['node_type'] -%>
<%# HeadNode -%>
<% when 'HeadNode' -%>
<% if node['cluster']['scheduler'] == 'slurm' -%>
[program:clustermgtd]
command = <%= node_virtualenv_path %>/bin/clustermgtd
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"<% if @region.start_with?('us-iso') -%>,AWS_CA_BUNDLE="<%= @aws_ca_bundle %>"<% end -%>
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/clustermgtd
stdout_logfile_maxbytes = 0
<% end -%>
<% unless node['cluster']['scheduler'] == 'awsbatch' -%>
[program:clusterstatusmgtd]
command = <%= cookbook_virtualenv_path %>/bin/python /opt/parallelcluster/scripts/clusterstatusmgtd.py
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"<% if @region.start_with?('us-iso') -%>,AWS_CA_BUNDLE="<%= @aws_ca_bundle %>"<% end -%>
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/clusterstatusmgtd
stdout_logfile_maxbytes = 0
<% end -%>
<% if @dcv_configured -%>
[program:pcluster_dcv_authenticator]
command = <%= @dcv_auth_virtualenv_path %>/bin/python <%= @dcv_auth_user_home %>/pcluster_dcv_authenticator.py
      --port <%= Integer(@dcv_port) + 1 %>
      --certificate <%= @dcv_auth_certificate %>
      --key <%= @dcv_auth_private_key %>
user = <%= @dcv_auth_user %>
environment = HOME="<%= @dcv_auth_user_home %>",USER="<%= @dcv_auth_user %>"<% if @region.start_with?('us-iso') -%>,AWS_CA_BUNDLE="<%= @aws_ca_bundle %>"<% end -%>
<% end -%>

<%# ComputeFleet -%>
<% when 'ComputeFleet' -%>
<% if node['cluster']['scheduler'] == 'slurm' -%>
[program:computemgtd]
command = <%= node_virtualenv_path %>/bin/computemgtd
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"<% if @region.start_with?('us-iso') -%>,AWS_CA_BUNDLE="<%= @aws_ca_bundle %>"<% end -%>
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/computemgtd
stdout_logfile_maxbytes = 0
<% end -%>

<%# LoginNode -%>
<% when 'LoginNode' -%>
[program:loginmgtd]
command = <%= node['cluster']['shared_dir_login_nodes'] %>/loginmgtd.sh
user = <%= node['cluster']['cluster_admin_user'] %>
environment = HOME="/home/<%= node['cluster']['cluster_admin_user'] %>",USER="<%= node['cluster']['cluster_admin_user'] %>"
autorestart = unexpected
exitcodes = 0
redirect_stderr = true
stdout_logfile = /var/log/parallelcluster/loginmgtd.log
stdout_logfile_maxbytes = 1MB
<% end -%>
