# templates/ubuntu/cfncluster_supervisord.conf.erb
# Generated by Chef for cfncluster <%= node['cfncluster']['cfn_node_type'] -%>
# Local modifications could be be overwritten.
<% case node['cfncluster']['cfn_node_type'] -%>
<%# MasterServer -%>
<% when 'MasterServer' -%>
[program:sqswatcher]
command = /usr/local/bin/sqswatcher
redirect_stderr = true
stdout_logfile = /var/log/sqswatcher

<%# ComputeFleet -%>
<% when 'ComputeFleet' -%>
[program:nodewatcher]
command = /usr/local/bin/nodewatcher
redirect_stderr = true
stdout_logfile = /var/log/nodewatcher
<% end -%>